# alarm.cron
* * * * * sar -u 60 3 | awk -f /home/user/day2/cpu_usage.awk | xargs -n 2 /home/user/day2/alarm.sh
* * * * * sar -r 60 3 | awk -f /home/user/day2/memory_usage.awk | xargs -n 2 /home/user/day2/alarm.sh
* * * * * sar -n DEV 1 60 | awk -f /home/user/day2/network_usage.awk | xargs -n 2 /home/user/day2/alarm.sh
